# target/linux/x86/64/base-files/etc/hotplug.d/net/10-mlxsw-spectrum-labels
#!/bin/sh

if [[ "$ACTION" = "add" && "$DEVICEDRIVER" = "mlxsw_spectrum*" ]]; then
	read -r label < "/sys/class/net/$DEVICENAME/phys_port_name"
	ip link set "$DEVICENAME" name "sw$label"
fi
